# launchd/com.local.vortex.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.local.vortex</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<!-- 把下面整行换成你的 vortex 绝对路径；sleep 可改短/加长 -->
		<string>sleep 15; REPLACE_WITH_ABSOLUTE_PATH_TO_VORTEX/start_vortex.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
